# platform/systemd/openclaw-browserlab.service
[Unit]
Description=OpenClaw Browser Lab
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=__REPO_ROOT__
Environment=OPENCLAW_STATE_DIR=__STATE_DIR__
ExecStart=/bin/bash -lc 'docker compose -f __REPO_ROOT__/platform/compose/docker-compose.browser-lab.yaml up -d'
ExecStop=/bin/bash -lc 'docker compose -f __REPO_ROOT__/platform/compose/docker-compose.browser-lab.yaml down'

[Install]
WantedBy=default.target
