# platform/systemd/openclaw-browserlab.service
[Unit]
Description=OpenClaw Browser Lab
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=__REPO_ROOT__
Environment=OPENCLAW_STATE_DIR=__STATE_DIR__
Environment=PATH=%h/.config/varlock/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
ExecStart=__PYTHON_EXECUTABLE__ -m clawops ops --asset-root __REPO_ROOT__ browser-lab up
ExecStop=__PYTHON_EXECUTABLE__ -m clawops ops --asset-root __REPO_ROOT__ browser-lab down

[Install]
WantedBy=default.target
