# platform/systemd/openclaw-sidecars.service
[Unit]
Description=OpenClaw Auxiliary Sidecars
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=__REPO_ROOT__
Environment=OPENCLAW_STATE_DIR=__STATE_DIR__
ExecStart=/bin/bash -lc '__REPO_ROOT__/scripts/ops/launch_sidecars_with_varlock.sh'
ExecStop=/bin/bash -lc '__REPO_ROOT__/scripts/ops/stop_sidecars.sh'

[Install]
WantedBy=default.target
