# platform/systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__REPO_ROOT__
Environment=OPENCLAW_STATE_DIR=__STATE_DIR__
ExecStart=/bin/bash -lc '__REPO_ROOT__/scripts/ops/launch_gateway_with_varlock.sh'
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=default.target
